# serveraccess.conf
#this one goes to /etc/init/
description "server access"
start on runlevel [23]
start on started tty1
respawn
respawn limit 10 5
exec sudo -u www-data /var/www/ServerAccess/runserver.sh